# dalybms_waveshare/rootfs/etc/services.d/dalybms/run
#!/usr/bin/with-contenv bashio

export DEVICE=$(bashio::config 'device')
export DEVICE_ID=$(bashio::config 'device_id')
export CELL_COUNT=$(bashio::config 'cells_in_series')
export MQTT_SERVER=$(bashio::config 'mqtt_server')
export MQTT_USER=$(bashio::config 'mqtt_user')
export MQTT_PASS=$(bashio::config 'mqtt_pass')
export MQTT_CLIENT_ID=$(bashio::config 'mqtt_client_id')
export MQTT_DISCOVERY_PREFIX=$(bashio::config 'mqtt_discovery_prefix')
export POLL_INTERVAL_SECONDS=$(bashio::config 'poll_interval_seconds')

echo "=== Daly Waveshare service starting ==="
echo "DEVICE=$DEVICE"
echo "DEVICE_ID=$DEVICE_ID"
echo "CELL_COUNT=$CELL_COUNT"
echo "MQTT_SERVER=$MQTT_SERVER"
echo "MQTT_CLIENT_ID=$MQTT_CLIENT_ID"
echo "MQTT_DISCOVERY_PREFIX=$MQTT_DISCOVERY_PREFIX"
echo "POLL_INTERVAL_SECONDS=$POLL_INTERVAL_SECONDS"

exec python3 -u /monitor.py
